# lib/systemd/system/tor.service.d/40_qubes.conf
## This file is part of Qubes+Whonix.
## Copyright (C) 2015 Jason Mehring <nrgaway@gmail.com>
## License: GPL-2+

[Unit]
ConditionPathExists=/var/run/qubes-service/whonix-gateway
ConditionPathExists=!/var/run/qubes-service/whonix-firewall-failed

## https://github.com/marmarek/qubes-core-agent-linux/blob/master/vm-systemd/tor.service.d/30_qubes.conf
After=qubes-sysinit.service
