# lib/systemd/system/privoxy.service.d/40_qubes.conf
## Copyright (C) 2020 - 2023 ENCRYPTED SUPPORT LP <adrelanos@whonix.org>
## See the file COPYING for copying conditions.

[Unit]
## https://forums.whonix.org/t/installation-and-fix-of-i2p-inside-whonix-workstation-by-default/8610
ConditionPathExists=!/run/qubes/this-is-templatevm
