# lib/systemd/system/tor@default.service.d/40_qubes.conf
## This file is part of Qubes+Whonix.
## Copyright (C) 2015 Jason Mehring <nrgaway@gmail.com>
## License: GPL-2+

[Unit]
ConditionPathExists=/run/qubes-service/whonix-gateway

## https://github.com/marmarek/qubes-core-agent-linux/pull/74
After=qubes-sysinit.service
