# lib/systemd/system/sdwdate.service.d/40_qubes.conf
## Copyright (C) 2015 Jason Mehring <nrgaway@gmail.com>
## License: GPL-2+

[Unit]
ConditionPathExists=!/run/qubes/this-is-templatevm
